# scripts/mozartd.service
[Unit]
Description=Mozart AI Compose Daemon
After=network.target
Documentation=https://github.com/emzi/mozart-ai-compose

[Service]
Type=notify
ExecStart=/usr/local/bin/mozartd start --foreground
ExecStop=/usr/local/bin/mozartd stop
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
WatchdogSec=60
User=%i
Group=%i

# Environment
Environment=PYTHONUNBUFFERED=1

# Resource limits
MemoryMax=8G
TasksMax=100

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
